# crontab/nsw-rfs-majorincidents-geojson
# update every 5 minutes
*/5 *  *   *   *  root  nsw-rfs-majorincidents-geojson /srv/www/nsw-rfs-majorincidents.geojson
